# ansible/roles/zeek/templates/promiscuous_nic.service.j2
# this file was generated. edits will be lost
#
[Unit]
Description=Makes an interface run in promiscuous mode at boot
After=network.target

[Service]
Type=oneshot
ExecStart=/sbin/ip link set dev {{ zeek_interface }} promisc on
TimeoutStartSec=0
RemainAfterExit=yes

[Install]
WantedBy=default.target
